# etc/supervisor.d/lighttpd.ini
[program:lighttpd]
command=bash -c "sleep 5 && /usr/sbin/lighttpd -D -f /etc/lighttpd/lighttpd.conf"
redirect_stderr=true