# etc/supervisor.d/logpipe.ini
[program:logpipe]
command=/bin/cat /tmp/lighttpd.logpipe
user=lighttpd
redirect_stderr=true