# etc/supervisor.d/smtp-sink.ini
[program:smtp-sink]
command=/usr/bin/smtp-sink -c -d "/home/smtp/MailDir/new/%%M" -u smtp 0.0.0.0:25 10
redirect_stderr=true
